# Week12/CesarOmar/3/space-check.service
[Unit]
After=mysql.service

[Service]
ExecStart=/usr/local/bin/space-check.sh

[Install]
WantedBy=default.target
