# Week12/ParshotanSeenanan/programTwo.service
[Unit]
Description=This program helps me with the Python Challenge

[Service]
ExecStart=/opt/programTwo.py

[Install]
WantedBy=multi-user.target
