# Week12/ParshotanSeenanan/programOne.service
[Unit]
Description=Counts to 50 and 100

[Service]
ExecStart=/opt/programOne.py

[Install]
WantedBy=multi-user.target
