# Week12/ParshotanSeenanan/programThree.service
[Unit]
Description=Generates a random complex password

[Service]
ExecStart=/opt/programThree.py

[Install]
WantedBy=multi-user.target

