# Week12/StevenSurujballi/Program1/hello.service
[Unit]
Description=My first service

[Service]
ExecStart=/opt/hello.py

[Install]
WantedBy=multi-user.target
