# Week12/MarkDecello/bash_increment.service
[Unit]
Description=Bash Service
After=increment.service

[Service]
ExecStart=/opt/bash_increment.sh

[Install]
WantedBy=multi-user.target
