# Week12/MarkDecello/hello.service
[Unit]
Description=My first service
After=bash_increment.service

[Service]
ExecStart=/opt/hello.py

[Install]
WantedBy=multi-user.target
