# Week12/MarkDecello/increment.service
[Unit]
Description=Python Service

[Service]
ExecStart=/opt/increment.py

[Install]
WantedBy=multi-user.target
