# Week12/BrunoBrito/my_first_service.service
[Unit]
Description=My first service

[Service]
ExecStart=/opt/script1.sh

[Install]
WantedBy=multi-user.target
