# Week12/BrunoBrito/my_second_service.service
[Unit]
Description=My second service

[Service]
ExecStart=/opt/script2.sh

[Install]
WantedBy=multi-user.target
