# Week12/BrunoBrito/my_third_service.service
[Unit]
Description=My third service

[Service]
ExecStart=/opt/script3.sh

[Install]
WantedBy=multi-user.target
