# Week12/Luis Collymoore/Service1/km.service
Unit]
Description=First Service

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/km.sh

[Install]
WantedBy=multi-user.target
