# Week12/StevenSurujballi/Program3/pythonloop.service
[Unit]
Description=Python program generates random characters

[Service]
ExecStart=/opt/pythonloop.py

[Install]
WantedBy=multi-user.target
