# Week12/Richardhw12/Feliz.service
[Unit]
Description=This is my second service HA! Take that

[Service]
ExecStart=/opt/Feliz.py

[Install]
WantedBy=multi-user.target
