# Week12/Richardhw12/hello.service
[Unit]
Description=Richard First Service

[Service]
ExecStart=/opt/hello.py

[Install]
WantedBy=multi-user.target
