# Week12/Richardhw12/loud.service
[Unit]
Description=LA Tercera Gano Rich!!! 

[Service]
ExecStart=/opt/loud.py

[Install]
WantedBy=multi-user.target
