# Week12/Michael_zambrano/service1/my_first_service.service
[Unit]
Description=First Sys Service 

[Service]
ExecStart=/bin/bash /usr/bin/first_service.sh
Type=simple

[Install]
WantedBy=multi-user.target
