# Week12/ElizabethRappleyea/fetcher.service
[Unit]
Description=My Personal Updater

[Service]
ExecStart=/opt/initfetch.sh

[Install]
WantedBy=multi-user.target
