# Week12/ElizabethRappleyea/funfacts.service
[Unit]
Description="first homework service"

[Service]
ExecStart=/opt/funfacts.py

[Install]
WantedBy=multi-user.target
