# Week12/ElizabethRappleyea/motfd.service
[Unit]
Description=A shell serivce

[Service]
ExecStart=/opt/motfd.sh

[Install]
WantedBy=multi-user.target
