# Week12/SebastianTapia/hi.service
[Unit]
Description=My first service

[Service]
ExecStart=/opt/hi.py

[Install]
WantedBy=multi-user.target 
