# Week12/SebastianTapia/myfirst.service
[Unit]
Description=primer service

[Service]
ExecStart=/opt/first.sh

[Install]
WantedBy=multi-user.target
