# Week12/SebastianTapia/third.service
[Unit]
Description=My last service

[Service]
ExecStart=/opt/thescript3.sh

[Install]
WantedBy=multi-user.target
