# Week12/MackensieAlvarez/p2.service
[Unit]
Description=This program does the area of a triangle 

[Service]
ExecStart=/opt/p2.py

[Install]
WantedBy=multi-user.target
