# Week12/MackensieAlvarez/p3.service
[Unit]
Description=This program sorts a list

[Service]
ExecStart=/opt/p3.py

[Install]
WantedBy=multi-user.target
