# Week12/Luis Collymoore/Service3/while.service
[Unit]
Description=While loop that counts up to 5

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/while.sh

[Install]
WantedBy=multi-user.target
